# src/rawdata-downloader/src/system-tools/etc/udev/rules.d/90-rawdata-downloader.rules
KERNEL=="sd*", RUN+="__PREFIX__/usr/sbin/scsi-disk-udev-eventhandler"

